# ansible/roles/frigate_vm/files/krt-container-watchdog.service
[Unit]
Description=KRT allowlisted container watchdog
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/krt-container-watchdog.sh
TimeoutStartSec=240
RuntimeDirectory=krt-container-watchdog
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/krt-container-watchdog
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX
